# scripts/com.hms.alertas-criticas.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hms.alertas-criticas</string>

    <key>ProgramArguments</key>
    <array>
        <string>REEMPLAZAR_RUTA_PROYECTO/.venv/bin/python3</string>
        <string>REEMPLAZAR_RUTA_PROYECTO/scripts/alertas_criticas.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>REEMPLAZAR_RUTA_PROYECTO</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- IMPORTANTE: los logs deben estar FUERA de ~/Documents/ porque
         macOS TCC bloquea a launchd para escribir ahí — causa exit
         EX_CONFIG (78) antes de que Python arranque. -->
    <key>StandardOutPath</key>
    <string>REEMPLAZAR_HOME/Library/Logs/HMS/launchd-criticas-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>REEMPLAZAR_HOME/Library/Logs/HMS/launchd-criticas-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
