# scripts/com.hms.alertas-semanales.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hms.alertas-semanales</string>

    <!-- Cron semanal — lunes 7:30 AM. Manda a cada cliente un email
         con el pronóstico de los próximos 7 días + alertas previstas
         + acciones generales. Es preventivo: ayuda al cliente a
         planificar la semana antes de que empiecen las alertas
         puntuales. -->

    <key>ProgramArguments</key>
    <array>
        <string>REEMPLAZAR_RUTA_PROYECTO/.venv/bin/python3</string>
        <string>REEMPLAZAR_RUTA_PROYECTO/scripts/alertas_semanales.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>REEMPLAZAR_RUTA_PROYECTO</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- IMPORTANTE: los logs deben estar FUERA de ~/Documents/ porque
         macOS TCC bloquea a launchd para escribir ahí — causa exit
         EX_CONFIG (78) antes de que Python arranque. -->
    <key>StandardOutPath</key>
    <string>REEMPLAZAR_HOME/Library/Logs/HMS/launchd-semanal-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>REEMPLAZAR_HOME/Library/Logs/HMS/launchd-semanal-stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
