# scripts/com.hms.pedido-carga.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hms.pedido-carga</string>

    <!-- Cron CADA 15 MINUTOS — el script evalúa qué lotes están en
         ventana de envío (±10 min de su hora_comida_1 / hora_comida_2,
         o del recordatorio +60 min si no cargaron). Soporta múltiples
         comidas/día configurables por lote.
         Para silocomedero solo manda el día que toca recargar. -->

    <key>ProgramArguments</key>
    <array>
        <string>REEMPLAZAR_RUTA_PROYECTO/.venv/bin/python3</string>
        <string>REEMPLAZAR_RUTA_PROYECTO/scripts/whatsapp_pedido_carga.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>REEMPLAZAR_RUTA_PROYECTO</string>

    <!-- Cada 15 minutos = 900 segundos -->
    <key>StartInterval</key>
    <integer>900</integer>

    <!-- IMPORTANTE: los logs deben estar FUERA de ~/Documents/ porque
         macOS TCC bloquea a launchd para escribir ahí — causa exit
         EX_CONFIG (78) antes de que Python arranque. -->
    <key>StandardOutPath</key>
    <string>REEMPLAZAR_HOME/Library/Logs/HMS/launchd-pedido-carga-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>REEMPLAZAR_HOME/Library/Logs/HMS/launchd-pedido-carga-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
